# etc/rc.d/rc.local
#!/bin/sh
#
# /etc/rc.d/rc.local:  Local system initialization script.
#
# Put any local startup commands in here.  Also, if you have
# anything that needs to be run at shutdown time you can
# make an /etc/rc.d/rc.local_shutdown script and put those
# commands in there.

export ARCH=`/bin/arch`
/usr/bin/loadkeys ~jimmy/etc/br-abnt.map.gz

for service in firewall mongodb postgresql tor; do
    test -x /etc/rc.d/rc.$service || continue
    (/etc/rc.d/rc.$service start) >/dev/null 2>&1 &
done
